# codex-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template launchd service definition for codex-proxy.

  Prefer:
    codex-proxy install

  The install command generates ~/Library/LaunchAgents/com.local.codex-proxy.plist
  with the current user's HOME, binary path, and log path.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.codex-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>__CODEX_PROXY_BIN__</string>
        <string>serve</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__CODEX_PROXY_LOG__</string>
    <key>StandardErrorPath</key>
    <string>__CODEX_PROXY_LOG__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__CODEX_PROXY_HOME__</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
