# codex-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd service definition for codex-proxy.

  Install:
    cp codex-proxy.plist ~/Library/LaunchAgents/com.local.codex-proxy.plist
    launchctl load ~/Library/LaunchAgents/com.local.codex-proxy.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.local.codex-proxy.plist
    rm ~/Library/LaunchAgents/com.local.codex-proxy.plist

  Logs:
    tail -f /usr/local/var/log/codex-proxy.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.codex-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/codex-proxy</string>
        <string>serve</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/codex-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/codex-proxy.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/peter</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
